# digitalocean-ubuntu-oscar/etc/supervisor/gunicorn.conf
[program:gunicorn]
command=%(gunicorn_path)s %(project)s.wsgi:application -c %(project_dir)s/gunicorn.conf.py
directory=%(project)s
user=nobody
autorestart=true
redirect_stderr=true