# digitalocean-ubuntu-searchengines/config/supervisor/solr.conf
[program:solr]
directory=%(solr_path)s
command=java -Djava.util.logging.config.file=logging.properties -jar start.jar