# digitalocean-ubuntu-searchengines/config/supervisor/sphinx.conf
[program:sphinxsearch]
command=searchd --config /etc/sphinxsearch/sphinx.conf
autostart=true
autorestart=true
stderr_logfile=/var/log/sphinx.err.log
stdout_logfile=/var/log/sphinx.out.log